# SetUpScheduler/crontab
# START CRON JOB
* * * * * python task1/dags/task1.py 
# END CRON JOB